# root/etc/services.d/monitoring/run
#!/usr/bin/with-contenv bash

cd /config

exec \
  s6-setuidgid grfd \
    /usr/bin/java -Dlog4j2.formatMsgNoLookups=true -jar \
      /app/monitoring/monitoring.jar
