# root/etc/services.d/monitoring/run
#!/usr/bin/with-contenv bash

cd /app/monitoring

exec \
  s6-setuidgid grfd \
    /usr/bin/java -jar \
      monitoring.jar